# mount.ondemand/files/mount-ondemand.init
#!/bin/sh /etc/rc.common

STOP=95

stop () {
	killall mount.ondemand
}

